# mkdocs/rootfs/etc/s6-overlay/s6-rc.d/mkdocs/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: MkDocs
# Runs the MkDocs
# ==============================================================================
echo "Starting MkDocs Server"

# Start mkdocs server
. /opt/venv/bin/activate 
mkdocs serve -a 0.0.0.0:8000 -f /config/mkdocs.yml &

# Wait for mkdocs temporary directory is created
sleep 1

# Start nginx to serve page
rm /data/www -fr
mkdir -p /data/
ln -s /tmp/mkdocs_* /data/www
chmod 777 /data/www
chmod 777 /tmp/mkdocs_*
nginx -g "daemon off; error_log /dev/stdout debug;"
