# mkdocs/rootfs/etc/s6-overlay/s6-rc.d/mkdocs/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: MkDocs
# Runs the MkDocs
# ==============================================================================
echo "Starting MkDocs Server"

. /opt/venv/bin/activate 
mkdocs serve  -a 0.0.0.0:8000 -f /config/mkdocs.yml &
sleep 1
cd /tmp/mkdocs_*
python -m http.server 9000 --protocol HTTP/1.1